# deploy/on-prem/production/openbrain-api.service
[Unit]
Description=OpenBrain API server (production system service)
Documentation=https://github.com/TurboAnbeli/OpenBrain
Wants=network-online.target postgresql.service ollama.service openbrain-embedder.service
After=network-online.target postgresql.service ollama.service openbrain-embedder.service
StartLimitIntervalSec=5min
StartLimitBurst=5

[Service]
Type=simple
User=openbrain
Group=openbrain
WorkingDirectory=/opt/openbrain/current
EnvironmentFile=/etc/openbrain/openbrain.env
Environment=NODE_ENV=production
Environment=XDG_CACHE_HOME=/var/cache/openbrain
Environment=TRANSFORMERS_CACHE=/var/cache/openbrain/transformers
ExecStartPre=/usr/local/lib/openbrain/openbrain-embedder-healthcheck.sh
ExecStart=/usr/bin/node /opt/openbrain/current/dist/index.js
Restart=on-failure
RestartSec=5s
TimeoutStartSec=90s
TimeoutStopSec=30s
KillSignal=SIGTERM

StateDirectory=openbrain
CacheDirectory=openbrain
LogsDirectory=openbrain
RuntimeDirectory=openbrain
UMask=0077

MemoryHigh=1536M
MemoryMax=2G
CPUQuota=300%
TasksMax=256

NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/openbrain /var/cache/openbrain /var/log/openbrain /run/openbrain
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
