# deploy/on-prem/systemd/consolidation-worker.timer
[Unit]
Description=Run OpenBrain consolidation worker every 15 minutes (single cycle)
Documentation=https://github.com/TurboAnbeli/OpenBrain
After=openbrain-api.service

[Timer]
OnCalendar=*:0/15
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
