# deploy/on-prem/systemd/consolidation-worker-oneshot.service
[Unit]
Description=OpenBrain consolidation worker — single cycle (invoked by timer)
Documentation=https://github.com/TurboAnbeli/OpenBrain
After=openbrain-api.service
Wants=openbrain-api.service

[Service]
Type=oneshot
WorkingDirectory=%h/workspace/openbrain
ExecStart=%h/workspace/openbrain/deploy/on-prem/systemd/start-consolidation-worker-cycle.sh
TimeoutStartSec=300
# Environment overrides (defaults are in .env):
# OPENBRAIN_LLM_CONSOLIDATION_ENDPOINT   — dedicated LLM endpoint for consolidation
# OPENBRAIN_LLM_CONSOLIDATION_MODEL      — LLM model for consolidation (default gemma-4-E4B-it)
# OPENBRAIN_SYNTHESIS_MODEL              — legacy fallback for OPENBRAIN_LLM_CONSOLIDATION_MODEL
# OLLAMA_ENDPOINT                         — fallback Ollama endpoint (default http://127.0.0.1:11434)
