# deploy/on-prem/systemd/consolidation-worker.service
[Unit]
Description=OpenBrain consolidation worker daemon (auto-discovers and consolidates thoughts)
Documentation=https://github.com/TurboAnbeli/OpenBrain
After=openbrain-api.service
Wants=openbrain-api.service

[Service]
Type=simple
WorkingDirectory=%h/workspace/openbrain
ExecStart=%h/workspace/openbrain/deploy/on-prem/systemd/start-consolidation-worker.sh
Restart=on-failure
RestartSec=30
# Give the worker time to connect to DB on first start.
TimeoutStartSec=120
# Environment overrides (defaults are in .env):
# CONSOLIDATION_INTERVAL_MS — poll interval in ms (default 900000 = 15 min)
# OPENBRAIN_SYNTHESIS_MODEL — LLM for synthesis (default qwen3:1.7b)
# OLLAMA_ENDPOINT — Ollama endpoint (default http://127.0.0.1:11434)

[Install]
WantedBy=default.target
